# angstrom/systemd/isoblue-interfaces.service
# /etc/systemd/system/isoblue-interfaces.service

[Unit]
Description=Bring CAN Interfaces Up
After=sysinit.target

[Service]
Type=oneshot
ExecStartPre=/sbin/ifconfig ib_eng down
ExecStartPre=/sbin/ifconfig ib_imp down
ExecStart=/bin/ip link set ib_eng type can bitrate 250000
ExecStart=/bin/ip link set ib_imp type can bitrate 250000
ExecStartPost=/sbin/ifconfig ib_eng up
ExecStartPost=/sbin/ifconfig ib_imp up
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
